# ops/systemd/bi-api.service.template
[Unit]
Description=BI API Service
After=network.target bi-redis.service
Wants=bi-redis.service

[Service]
Type=simple
User=root
WorkingDirectory=__ROOT_DIR__
Environment="NODE_ENV=production"
ExecStart=/usr/bin/env bash -lc 'cd __ROOT_DIR__ && source ./scripts/_test_server_env.sh && exec npm --workspace @flcbi/api run start'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
